# root/etc/services.d/darkhttpd/run
#!/usr/bin/with-contenv bash

#export HOME="~app"
#exec s6-setuidgid root "/conf/start_darkhttpd.sh"
exec /conf/start_darkhttpd.sh
